# examples/odoo-knowledge-rs.service
[Unit]
Description=Odoo Knowledge RS MCP Server
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
Environment=ODOO_KNOWLEDGE_CONFIG=<CONFIG_DIR>/production.toml
WorkingDirectory=<DATA_DIR>
ExecStart=<INSTALL_BIN>/odoo-knowledge serve
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=<DATA_DIR>

[Install]
WantedBy=multi-user.target
